# docker/mongo/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=mongod --bind_ip_all

[program:daemon]
command=python3 daemon.py -a restart
